# python/pi-firmware/boot/autopilot-pi.service
[Unit]
Description=Autopilot Pi Central Server (PM2)
Documentation=https://github.com/Threedreamz/master-weaver-autopilot-starter
After=network-online.target avahi-daemon.service nginx.service
Wants=network-online.target
Requires=nginx.service

[Service]
Type=forking
User=pi
Group=pi
WorkingDirectory=/opt/autopilot
ExecStart=/usr/bin/pm2 start /opt/autopilot/ecosystem.config.js
ExecReload=/usr/bin/pm2 reload all
ExecStop=/usr/bin/pm2 stop all
Restart=always
RestartSec=5

# Environment
Environment=NODE_ENV=production
Environment=PM2_HOME=/home/pi/.pm2

# Resource limits
MemoryMax=1G
CPUQuota=90%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=autopilot-pi

[Install]
WantedBy=multi-user.target
